# tmp/droid-config.files
%defattr(-,root,root,-)
/usr/lib/sysctl.d/ipv6-settings.conf
/usr/lib/sysctl.d/ipv4-settings.conf
/usr/lib/systemd/user/jolla-startupwizard-pre-user-session.service.d/50-compositor.conf
/usr/lib/systemd/user/lipstick.service.d/50-compositor.conf
/usr/lib/systemd/user/jolla-actdead-charging.service.d/50-compositor.conf
/usr/lib/systemd/system/timers.target.wants/droid-reclaim-memory.timer
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/droid-late-start.service
/usr/lib/systemd/system/droid-reclaim-memory.service
/usr/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-late-start.service
/usr/lib/systemd/system/actdead.target.wants/droid-hal-init.service
/usr/lib/systemd/system/droid-reclaim-memory.timer
/usr/lib/tmpfiles.d/droid-hal.conf
/usr/bin/droid/droid-init-done.sh
/usr/bin/droid/get_partition.sh
/usr/bin/droid/droid-load-firmware.sh
/usr/bin/droid/droid-reclaim-memory.sh
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/droid-hal-shutdown.sh
/lib/udev/rules.d/997-droid-drm.rules
/lib/udev/rules.d/998-droid-system.rules
/etc/lowi.conf
/etc/xdg/QtProject/QPlatformTheme.conf
/etc/xdg/QtProject/Sensors.conf
/etc/xdg/QtProject/QtQuick2.conf
/etc/obexd/root
/etc/obexd/opp.conf
/etc/obexd/excludes/bluetooth:pcsuite
/etc/obexd/noplugins/ftp
/etc/obexd/noplugins/syncevolution
/etc/obexd/noplugins/irmc
/etc/obexd/noplugins/mas
/etc/obexd/pbap.conf
/etc/media_codecs_performance.xml
/etc/selinux/config
/etc/izat.conf
/etc/acdbdata
/etc/sensorfw/primaryuse.conf
/etc/usb-moded/dyn-modes/mtp_mode-droid.ini
/etc/media_codecs_google_telephony.xml
/etc/gps.conf
/etc/sap.conf
/etc/mixer_paths.xml
/etc/modprobe.d/blacklist-evbug.conf
/etc/media_codecs.xml
/etc/ofono/binder.d/dual-sim.conf
/etc/ofono/binder.conf
/etc/media_codecs_ffmpeg.xml
/etc/audio_platform_info.xml
/etc/xtra_root_cert.pem
/etc/media_codecs_google_audio.xml
/etc/media_codecs_google_video.xml
/etc/xtwifi.conf
/etc/sec_config
/etc/media_profiles.xml
/var/lib/environment/compositor/droid-hal-device.conf
/var/lib/environment/nemo/60-multimedia.conf
/var/lib/environment/ofono/noplugin.conf
/usr/lib/systemd/system/multi-user.target.wants/wlan-module-load.service
/usr/lib/systemd/system/bluetooth.service.wants/droid-hcismd-up.service
/usr/lib/systemd/system/overlay-vendor.service
/usr/lib/systemd/system/system-etc-ld.config.28.txt.mount
/usr/lib/systemd/system/wlan-module-load.service
/usr/lib/systemd/system/droid-hcismd-up.service
/usr/lib/systemd/system/local-fs.target.wants/vendor.mount
/usr/lib/systemd/system/local-fs.target.wants/overlay-vendor.service
/usr/lib/systemd/system/local-fs.target.wants/system-etc-ld.config.28.txt.mount
/usr/bin/logcat
/usr/bin/droid/android-links.sh
/usr/bin/droid/droid-get-bt-address.sh
/usr/bin/droid/droid-hal-early-init.sh
/usr/bin/droid/mount_overlay_vendor.sh
/usr/share/csd/settings.d/hw-settings.ini
/usr/libexec/droid-hybris/system/etc/init/hwservicemanager.rc
/usr/libexec/droid-hybris/system/etc/init/permissions_fix.rc
/usr/libexec/droid-hybris/system/etc/init/disabled_services.rc
/usr/libexec/droid-hybris/system/etc/init/ecclist.rc
/usr/libexec/droid-hybris/system/etc/init/vndservicemanager.rc
/usr/libexec/droid-hybris/system/etc/ld.config.28.txt
/lib64/firmware
/lib/udev/rules.d/999-droid-fm.rules
/etc/audio_policy_configuration.xml
/etc/mce/60-input-policy-davinci.conf
/etc/mce/60-proximity-sensor-davinci.conf
/etc/mce/60-hybris-features-davinci.ini
/etc/gst-droid/gstdroidcodec.conf
/etc/connman/main-native.conf
/etc/selinux/minimum/contexts/dbus_contexts
/etc/selinux/minimum/contexts/files/file_contexts
/etc/selinux/minimum/policy/policy.30
/etc/selinux/config
/etc/acdbdata
/etc/sensorfw/sensord.conf.d/60-sensor-selection-overrides.conf
/etc/usb-moded/usb-moded-configfs.ini
/etc/usb-moded/dyn-modes/99-developer_mode-usb.ini
/etc/usb-moded/usb-moded.ini
/etc/ofono/ril_subscription.conf
/etc/ofono/ril_subscription.d/dual-sim.conf
/var/lib/environment/compositor/droid-hal-ginkgo.conf
/var/lib/environment/connman/override-main-config.conf
/var/lib/environment/nemo/99-qtscenegraph.conf
/var/lib/bluetooth/board-address
/usr/share/ssu/features.d/adaptation-community.ini
